# deploy/com.daily-news.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- 每周一 08:30 生成周趋势报告并推送飞书。安装：见 README「部署（launchd）」。
     注意把 PROJECT_DIR 占位符替换为绝对路径。 -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.daily-news.weekly</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>PROJECT_DIR/deploy/run.sh</string>
    <string>weekly</string>
    <string>--push</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>8</integer>
    <key>Minute</key><integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>PROJECT_DIR/logs/launchd.weekly.out.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_DIR/logs/launchd.weekly.err.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
